# systemd/promptworks-healthcheck.timer
[Unit]
Description=Run PromptWorks-OS health check daily

[Timer]
OnCalendar=*-*-* 09:00:00
Persistent=true

[Install]
WantedBy=timers.target
