# systemd/promptworks-update.service
[Unit]
Description=PromptWorks-OS Firmware-Works Update

[Service]
Type=oneshot
ExecStart=/bin/bash -lc '__PW_USER_HOME__/.local/bin/pw firmware apply'
