# systemd/promptworks-update.timer
[Unit]
Description=Run PromptWorks-OS firmware updates daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
